# smb-mount/rootfs/etc/services.d/cifs/run
#!/usr/bin/with-contenv bash

mount -t cifs -o $MOUNTOPTIONS $SERVERPATH $MOUNTPOINT

while true
do
  sleep 100
done