# rclone-mount/rootfs/etc/services.d/rclone-mount/finish
#!/usr/bin/with-contenv sh

s6-svwait -D /var/run/s6/services/rclone-mount/
#s6-svscanctl -t /var/run/s6/services
#s6-svc -wD /var/run/s6/services/rclone-mount/