# rclone-mount/rootfs/etc/services.d/rclone-mount/run
#!/usr/bin/with-contenv sh

/usr/sbin/rclone --config $ConfigPath mount $RemotePath $MountPoint $MountCommands